# gcr/gcr-ssh-agent.service.in
[Unit]
Description=GCR ssh-agent wrapper

Requires=gcr-ssh-agent.socket

[Service]
Type=simple
StandardError=journal
Environment=SSH_AUTH_SOCK=%t/gcr/ssh
ExecStart=@libexecdir@/gcr-ssh-agent --base-dir %t/gcr
Restart=on-failure

[Install]
Also=gcr-ssh-agent.socket
WantedBy=default.target
