# gcr/gcr-ssh-agent.socket.in
[Unit]
Description=GCR ssh-agent wrapper

[Socket]
Priority=6
Backlog=5
ListenStream=%t/gcr/ssh
ExecStartPost=-@systemctl@ --user set-environment SSH_AUTH_SOCK=%t/gcr/ssh
DirectoryMode=0700

[Install]
WantedBy=sockets.target
